# scripts/services/postgres-backup.service
# ── Creatrix PostgreSQL backup — systemd service unit ─────────────────────────
# canonical:   scripts/services/postgres-backup.service
# purpose:     Runs the backup script as a one-shot systemd service, triggered
#              by postgres-backup.timer. Not started directly.
#
# install:
#   sudo cp scripts/services/postgres-backup.service /etc/systemd/system/
#   sudo cp scripts/services/postgres-backup.timer   /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now postgres-backup.timer
#
# test (run immediately without waiting for timer):
#   sudo systemctl start postgres-backup.service
#   journalctl -u postgres-backup.service -f
#
# check schedule:
#   systemctl status postgres-backup.timer
# ──────────────────────────────────────────────────────────────────────────────

[Unit]
Description=Creatrix PostgreSQL backup
# Postgres must be running before the backup fires
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot

# Set this to the user that has read access to DATABASE_URL and write access
# to the backup directory. Adjust if you run Creatrix as a different user.
User=postgres

# Put DATABASE_URL and optionally CREATRIX_BACKUP_DIR / CREATRIX_BACKUP_KEEP
# in this file (one KEY=value per line, no export keyword, no quotes needed).
# Create it: sudo touch /etc/creatrix/backup.env && sudo chmod 600 /etc/creatrix/backup.env
EnvironmentFile=/etc/creatrix/backup.env

ExecStart=/bin/bash /path/to/creatrix/scripts/services/postgres-backup.sh

# Stream backup logs to the journal — view with:
#   journalctl -u postgres-backup.service
StandardOutput=journal
StandardError=journal
