# scripts/services/postgres-backup.timer
# ── Creatrix PostgreSQL backup — systemd timer unit ───────────────────────────
# canonical:   scripts/services/postgres-backup.timer
# purpose:     Schedules postgres-backup.service to run nightly at 02:00.
#              Adjust OnCalendar to change the schedule.
#
# install (run once):
#   sudo cp scripts/services/postgres-backup.service /etc/systemd/system/
#   sudo cp scripts/services/postgres-backup.timer   /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now postgres-backup.timer
#
# verify it's scheduled:
#   systemctl status postgres-backup.timer
#   systemctl list-timers postgres-backup.timer
#
# OnCalendar examples:
#   daily          → every day at midnight
#   *-*-* 02:00:00 → every day at 2 a.m.
#   weekly         → every Monday at midnight
#   *-*-1 03:00:00 → first of every month at 3 a.m.
# ──────────────────────────────────────────────────────────────────────────────

[Unit]
Description=Creatrix PostgreSQL nightly backup timer

[Timer]
# Run at 02:00 every night
OnCalendar=*-*-* 02:00:00

# If the machine was off at 02:00, run the backup within 10 minutes of next boot
Persistent=true

[Install]
WantedBy=timers.target
